# dist/restore.service.in
[Unit]
Description=Restore @project_name_dashes@ configuration

# When development gets far enough along that the driver creates its misc character device, uncomment these to wait
# for it.
# BindsTo=/dev/@project_name_dashes@
# After=/dev/@project_name_dashes@

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStartPre=@modprobe_path@ @project_name_underscores@
ExecStartPre=@udevadm_path@ settle
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/@project_name_dashes@-config --restore-config

ExecStop=@modprobe_path@ -r @project_name_underscores@

ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
